# com.bez.network-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.bez.network-router</string>

        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/network-router</string>
            <string>daemon</string>
            <string>-config</string>
            <string>/usr/local/etc/network-router/config.yaml</string>
        </array>

        <key>RunAtLoad</key>
        <true />

        <key>KeepAlive</key>
        <true />

        <!-- Set PATH so the daemon can find 'route', 'networksetup', etc. -->
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        </dict>

        <key>StandardOutPath</key>
        <string>/var/log/network-router.log</string>
        <key>StandardErrorPath</key>
        <string>/var/log/network-router.log</string>
    </dict>
</plist>